# application/luci-app-boostupnp/root/etc/uci-defaults/40_luci-boostupnp
#!/bin/sh

logger boostupnp "uci-default"
rm -f /tmp/luci-indexcache
exit 0
